# proctor/ansible/playbooks/files/reggie.conf
description "Reggie Registration portal"
author      "Rob Cameron"
console log

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on started networking
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/reggie"

    exec start-stop-daemon --start -d /reggie -c vagrant --exec  /usr/bin/python app.py runserver
end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "node.js has started!"
end script
